# hytale-panel/systemd/hytale-helper.service
# ============================================================
# Hytale Panel — Host Helper Service
# ============================================================
# This service runs on the HOST (not in Docker) and provides
# a Unix socket API for the panel to execute privileged
# operations: systemctl, tmux, file I/O, backups.
#
# It runs as a dedicated non-root hytale-helper user with a tight local-only systemd
# sandbox. Service control and journal reads are delegated through the
# narrow sudoers rules in hytale-helper.sudoers.
# Only the panel API container can connect to its Unix socket.
#
# Host socket path: /opt/hytale-panel/run/hytale-helper.sock
#
# Install:
#   sudo cp hytale-helper.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now hytale-helper.service
#
# Logs:
#   journalctl -u hytale-helper.service -f
# ============================================================

[Unit]
Description=Hytale Panel Host Helper
After=network.target
Documentation=https://github.com/your-repo/hytale-panel

[Service]
Type=simple
User=hytale-helper
Group=hytale-panel
SupplementaryGroups=hytale
WorkingDirectory=/opt/hytale-panel/helper
ExecStartPre=+/usr/bin/install -d -o hytale-helper -g hytale-panel -m 0770 /opt/hytale-panel/run
ExecStart=/usr/bin/node /opt/hytale-panel/helper/dist/index.js
Restart=always
RestartSec=5

# Environment file with HMAC secret and paths
EnvironmentFile=/opt/hytale-panel/helper/.env

# ─── Security hardening ───────────────────────────────────
#
# !!! READ BEFORE EDITING !!!
# The helper runs as a NON-ROOT user but MUST be able to invoke `sudo` for the
# narrow commands listed in /etc/sudoers.d/hytale-helper (systemctl on
# hytale-tmux.service and the validating journalctl wrapper). sudo is a
# setuid binary and CANNOT WORK when the kernel `no_new_privs` (NNP) bit is
# set on the process.
#
# Several systemd hardening directives IMPLICITLY set NoNewPrivileges=yes
# regardless of whether you also write `NoNewPrivileges=no`. Per
# `systemd.exec(5)` "If this option is used, NoNewPrivileges=yes is implied":
#   SystemCallFilter=, SystemCallArchitectures=, RestrictNamespaces=,
#   RestrictAddressFamilies= (when unprivileged), PrivateDevices=,
#   ProtectKernelTunables=, ProtectKernelModules=, ProtectKernelLogs=,
#   ProtectClock=, ProtectHostname=, MemoryDenyWriteExecute=,
#   RestrictRealtime=, RestrictSUIDSGID=, LockPersonality=, DynamicUser=
#
# Any of those will set NoNewPrivs:1 in /proc/$PID/status and break the
# Mods Manager Restart button (and every other helper sudoer call) with:
#   sudo: The "no new privileges" flag is set, which prevents sudo from
#         running as root.
#
# DO NOT add directives from that list back to this unit. The trade-off is
# accepted: the helper relies on filesystem confinement (ProtectSystem,
# ProtectHome, ReadWritePaths), namespace isolation (PrivateTmp), resource
# caps (MemoryMax, LimitNOFILE), and a non-root identity bounded by tight
# sudoers rules — instead of seccomp/BPF-class filters that require NNP.
#
# Verify after every edit:
#   sudo cp systemd/hytale-helper.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl restart hytale-helper.service
#   PID=$(systemctl show -p MainPID --value hytale-helper.service)
#   grep -E '^(NoNewPrivs|Uid)' /proc/$PID/status
# Expected: NoNewPrivs: 0   AND   Uid: <hytale-helper uid> ...
# scripts/doctor.sh enforces both at runtime.

# Explicit and load-bearing — do not flip without re-reading the comment above.
NoNewPrivileges=no

# AmbientCapabilities empty: helper inherits no capabilities. sudo brings its
# own setuid escalation path. This directive does NOT imply NNP.
AmbientCapabilities=

# Group-readable/writable files are intentional: backup and restore artifacts
# must remain accessible to the hytale game-server user through group perms.
UMask=0007

# Filesystem confinement — none of these imply NNP.
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
PrivateTmp=true
ReadWritePaths=/opt/hytale /opt/hytale-backups /opt/hytale-panel/run /opt/hytale-panel-data/mod-upload-staging

# Resource limits
LimitNOFILE=1024
MemoryMax=256M

[Install]
WantedBy=multi-user.target
